# source/server/conf/currant_task_any_supervisor_docker.conf
[program:currant_task_any]
command = /usr/bin/python wsgi.py -t any
directory=/var/lib/app/currant/source/server
user = 33
