# source/server/conf/currant_push_supervisor_docker.conf
[program:currant_push]
command = /usr/bin/python wsgi.py -u
directory=/var/lib/app/currant/source/server
user = 33
