# source/server/conf/currant_task_master_supervisor_docker.conf
[program:currant_task_master]
environment = LC_CTYPE="en_US.UTF-8"
command = /usr/bin/python wsgi.py -t master
directory = /var/lib/app/currant/source/server
user = 33
