# source/server/conf/currant_task_master_supervisor_docker.conf
[program:currant_task_master]
command = /usr/bin/python wsgi.py -t master
directory=/var/lib/app/currant/source/server
user = 33
