# source/server/conf/supervisor.docker.conf
[program:pyenv]
command = /usr/bin/docker run --rm --name pyenv -v /var/lib/app:/var/lib/app felixonmars/archlinux supervisord -c /var/lib/app/supervisor/supervisord.conf

